# etc/crontab
* */1 * * * /home/g/gorodin/gorod.in/manage.py count_article_comments_number >/dev/null